# root/etc/services.d/node/run
#!/usr/bin/with-contenv bash

echo "Starting Node actuator service"
node /app/actuator/server.js