# root/etc/services.d/nordvpnd/run
#!/usr/bin/with-contenv bash

base_dir="/vpn"
auth_file="$base_dir/auth"
config_file="$base_dir/config.ovpn"

random_sleep=$( shuf -i $MIN_RANDOM_SLEEP-$MAX_RANDOM_SLEEP -n 1 )s
echo "Delaying VPN startup for $random_sleep"
sleep $random_sleep

# Find a random server
/app/createvpnconfig.sh || exit 1

# On set up errors, the config file will not be created
if [[ ! -f "$config_file" ]]; then
	echo "VPN config file not present. Exiting."
	exit 1
fi

# Start the VPN service, and watch the stdout/stderr
# logs for authentication failures. If found, kill the container.
echo "Starting NordVPN"
exec openvpn --cd "$base_dir" \
	--config "$config_file" \
	--auth-user-pass "$auth_file" \
	--auth-nocache ${OPENVPN_OPTS} \
	2>&1