# root/etc/services.d/actuator/run
#!/usr/bin/with-contenv sh

echo "Starting Actuators"
exec node /app/node/actuator.js