# root/etc/services.d/actuator/run
#!/usr/bin/with-contenv sh

echo "Starting Node actuator service"
exec node /app/actuator/server.js