# root/etc/services.d/chrome/run
#!/usr/bin/with-contenv bash

# /usr/src/app
cd ${APP_DIR}

# Drop down to the unprivileged Chrome user
echo "Starting Chrome"
exec s6-setuidgid blessuser node ./build/index.js