# root/etc/services.d/crond/run
#!/usr/bin/env bash

echo "Starting cron"
exec /usr/sbin/cron -f
