# root/etc/services.d/crond/run
#!/usr/bin/env bash

exec /usr/sbin/cron -f
