# root/etc/services.d/openvpn/finish
#!/usr/bin/with-contenv sh

# Killing container due to VPN failure
s6-svscanctl -t /var/run/s6/services