# root/etc/services.d/randomizer/run
#!/usr/bin/env bash

echo "Starting randomizer"
while true ; do /app/randomizer/randomizer.sh & sleep 5 ; done